# templates/example.service.erb
[Unit]
Description=<%= scope['example::service_name'] %>
After=network.target

[Service]
Type=simple
ExecStart=<%= scope['example::install_dir'] %>/example > /dev/null 2>&1
ExecStop=<%= scope['example::install_dir'] %>/example > /dev/null 2>&1
